# cryptocloud-webhook.service
[Unit]
Description=CryptoCloud Payment Webhook Server
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/construct3store-bot
ExecStart=/usr/bin/python3 flask_webhook_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptocloud-webhook
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target